# launchd/com.gsg733jwb.cursor-sync-watch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gsg733jwb.cursor-sync-watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/cursor-dotfiles/scripts/cursor-sync-watch.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/.cursor/sync-watch.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.cursor/sync-watch.stderr.log</string>
</dict>
</plist>
